# qubes/systemd/qubes-netwatcher.service.d/corridor.conf
[Unit]
Require=corridor-init-forwarding.service
