# qubes/systemd/qubes-netwatcher.service.d/corridor.conf
[Unit]
Requires=corridor-init-forwarding.service
