# qubes/systemd/corridor-init-forwarding.service.d/qubes.conf
[Unit]
After=qubes-sysinit.service qubes-iptables.service
Before=qubes-firewall.service

[Service]
ExecStart=/bin/rm -f /var/run/qubes-service/qubes-firewall
