# qubes/systemd/qubes-iptables.service.d/corridor.conf
# break a cycle, see https://github.com/QubesOS/qubes-issues/issues/2209

[Unit]
DefaultDependencies=no
