# qubes/systemd/corridor-data.service.d/qubes.conf
[Unit]
After=qubes-mount-dirs.service
Wants=tor.service
