# systemd/corridor-init-logged.service.in
[Unit]
Description=corridor's logging
Requires=corridor-data.service
After=corridor-data.service

[Service]
ExecStart=@SBIN@/corridor-init-logged
Type=oneshot

[Install]
WantedBy=multi-user.target
