# systemd/corridor-init-snat.service.in
[Unit]
Description=corridor's source NAT
After=iptables.service

[Service]
ExecStart=SBIN/corridor-init-snat
ExecStop=SBIN/corridor-stop-snat
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
