# systemd/corridor-data.service.in
[Unit]
Description=corridor's relay list
After=tor.service

[Service]
ExecStart=SBIN/corridor-data
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
