# systemd/corridor-init-logged.service.in
[Unit]
Description=corridor's logging

[Service]
ExecStart=SBIN/corridor-init-logged
Type=oneshot

[Install]
WantedBy=multi-user.target
