# systemd/corridor-init-forwarding.service.in
[Unit]
Description=corridor's forwarding
Documentation=man:corridor(8)
DefaultDependencies=no
After=iptables.service systemd-sysctl.service
Before=network-pre.target
Wants=network-pre.target

[Service]
ExecStart=@SBIN@/corridor-init-forwarding
ExecStop=@SBIN@/corridor-stop-forwarding
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
RequiredBy=systemd-networkd.service
