# systemd/corridor-init-logged.service.in
[Unit]
Description=corridor's logging
Documentation=https://github.com/rustybird/corridor
Requires=corridor-data.service
After=corridor-data.service

[Service]
ExecStart=@SBIN@/corridor-init-logged
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
